# setup/golf_watcher_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Automatically monitors Round N Scores folders and updates the IMI Golf League Excel when new score files are saved.</Description>
    <Author>CORP+ehigh</Author>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <UserId>CORP\ehigh</UserId>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>CORP\ehigh</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>py</Command>
      <Arguments>-3 "C:\Users\ehigh\OneDrive - IMI Companies\Documents\Golf League\setup\watcher.py"</Arguments>
      <WorkingDirectory>C:\Users\ehigh\OneDrive - IMI Companies\Documents\Golf League</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
